# scripts/backup_db.cron
0 1 * * * /Users/randy/src/freedom_account/scripts/backup_db.sh
