# scripts/backup_db.cron
0 1 * * * /Users/randy/src/apps/freedom_account/scripts/backup_db.sh
