# scripts/backup_db.cron
0 2 * * * /Users/randy/src/apps/freedom_account/scripts/backup_db.sh
